# com.ai-supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-supervisor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lijunsheng/ai-supervisor/venv/bin/python</string>
        <string>/Users/lijunsheng/ai-supervisor/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/lijunsheng/ai-supervisor</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/lijunsheng/ai-supervisor/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/lijunsheng/ai-supervisor/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/lijunsheng/.nvm/versions/node/v22.22.0/bin:/Users/lijunsheng/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/lijunsheng</string>
    </dict>
</dict>
</plist>
